# monitoringProcess.timer
[Unit]
Description=Monitoring process timer

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
